# etc/rc.local
#!/bin/sh

# Start pigpiod at boot time; enable some gpios
# GPIO 18, 23 and 24, 25,7 and 8
# mask=0x1800000
mask=0x3840180
# -s 10 = read at 10us  (lower than default 5us)
/usr/bin/pigpiod -x $mask

# Start the sixad game controller driver; we need to wait for bt
# to be ready.
(sleep 20; 
hciconfig hci0 piscan
sixad --start & ) 
exit 0
